# configs/tuned/game-profile-monitor.path
[Unit]
Description=Watch for game processes

[Path]
DirectoryNotEmpty=/proc
TriggerLimitIntervalSec=1

[Install]
WantedBy=default.target